# launchd/com.devpulse.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.devpulse.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/darrencoxon/.bun/bin/bun</string>
        <string>--env-file</string>
        <string>../../.env</string>
        <string>src/index.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/darrencoxon/Documents/Codebases/current-projects/devPulse/apps/server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SERVER_PORT</key>
        <string>4000</string>
        <key>PATH</key>
        <string>/Users/darrencoxon/.bun/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/darrencoxon/Library/Logs/DevPulse/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/darrencoxon/Library/Logs/DevPulse/server.error.log</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
